# scripts/com.hotintelligence.daily-update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.hotintelligence.daily-update</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/huangxianfeng/Library/Application Support/HotIntelligence/run-daily-update.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/huangxianfeng/Documents/Codex/2026-05-19/new-chat</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>7</integer>
    <key>Minute</key>
    <integer>5</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/Users/huangxianfeng/Documents/Codex/2026-05-19/new-chat/logs/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/huangxianfeng/Documents/Codex/2026-05-19/new-chat/logs/launchd.err.log</string>
</dict>
</plist>
